# autostart/services/blueman.service
[Unit]
BindsTo=autostart.target
After=windowmanager.target

[Service]
ExecStart=/usr/bin/env blueman-applet
PassEnvironment=DISPLAY
