# autostart/services/dpms.service
[Unit]
BindsTo=autostart.target
After=windowmanager.target

[Service]
Type=oneshot
ExecStart=xset s off
ExecStart=xset dpms 0 0 300
