# autostart/services/dunst.service
[Unit]
BindsTo=autostart.target
After=windowmanager.target

[Service]
ExecStart=/usr/bin/env dunst -config %h/.config/dunstrc
Environment=DISPLAY=:0
