# autostart/services/dunst.service
[Unit]
BindsTo=autostart.target

[Service]
ExecStart=/usr/bin/dunst -config %h/.config/dunstrc
