# autostart/services/i3.service
[Service]
Type=exec
ExecStart=/usr/bin/env i3 --config %h/.i3/config
Environment=DISPLAY=:0

Restart=no
