# autostart/services/i3.service
[Service]
Type=exec
ExecStart=i3 --config %h/.i3/config

Restart=no
