# autostart/services/i3.service
[Service]
Type=exec
ExecStart=/usr/bin/env i3 --config %h/.i3/config
PassEnvironment=DISPLAY

Restart=no
