# autostart/services/nextcloud.service
[Unit]
BindsTo=autostart.target

ConditionEnvironment=MACHINE_HAS_NEXTCLOUD=true

[Service]
ExecStart=/usr/bin/nextcloud --background
