# autostart/services/nextcloud.service
[Unit]
BindsTo=autostart.target
After=windowmanager.target
After=gnome-keyring.service

ConditionEnvironment=MACHINE_HAS_NEXTCLOUD=true

[Service]
ExecStart=/usr/bin/env nextcloud --background
Environment=DISPLAY=:0
