# autostart/services/nm-applet.service
[Unit]
BindsTo=autostart.target
After=windowmanager.target

[Service]
ExecStart=/usr/bin/env nm-applet
Environment=DISPLAY=:0
