# autostart/services/pasystray.service
[Unit]
BindsTo=autostart.target
After=windowmanager.target

[Service]
ExecStart=/usr/bin/env pasystray
PassEnvironment=DISPLAY
