# autostart/services/pasystray.service
[Unit]
BindsTo=autostart.target

[Service]
ExecStart=/usr/bin/pasystray
