# autostart/services/redshift.service
[Unit]
BindsTo=autostart.target

[Service]
ExecStart=/usr/bin/redshift-gtk -c %h/.config/redshift.conf
