# autostart/services/restic.service
[Unit]
ConditionEnvironment=MACHINE_HAS_RESTIC_BACKUP=true

[Service]
Type=oneshot
ExecStart=%h/bin/restic-backup
RemainAfterExit=true
