# autostart/services/screencfg.service
[Unit]
BindsTo=autostart.target
After=windowmanager.target

[Service]
Type=simple
ExecStart=/usr/bin/screencfg watch --best
Restart=always
