# autostart/services/spotify.service
[Unit]
BindsTo=autostart.target

[Service]
ExecStart=/usr/bin/spotify
