# autostart/services/wallpaper.timer
[Unit]
BindsTo=autostart.target

[Timer]
OnActiveSec=0s
OnCalendar=hourly
Unit=wallpaper.service
RemainAfterElapse=true
