# autostart/services/windowmanager.target
[Unit]
Requires=i3.service
After=i3.service
