# autostart/services/yubikey-touch-detector.service
[Unit]
BindsTo=autostart.target

[Service]
ExecStart=/usr/bin/yubikey-touch-detector -libnotify
